# packaging/files/post-install.sh
#!/usr/bin/env bash

systemctl enable chirpstack-rest-api
systemctl start chirpstack-rest-api
